# Linux_os_config_sh/set_randM_linux_hostname&Ip/-scripts/rc.local
#!/bin/bash -e

# rc.local
# change file permission to root user only access:
# sudo chmod 700 start_fakeip

nohup /bin/bash /etc/udev/start_fakeip.sh

exit 0
